# roles/bastion/templates/etc/init/yubiauthd.conf
# {{ ansible_managed }}

description "yubiauthd"
author "Michael Sambol"

start on runlevel [2345]
stop on runlevel [!2345]

expect daemon
respawn

setuid root
setgid root

exec /usr/sbin/yubiauthd
