# roles/apt-mirror/templates/etc/cron.d/debmirror
# {{ ansible_managed }}
#
# Regular cron jobs for all packages that use debmirror syncing
#
0 6     * * *   apt-mirror      {{ apt_mirror.path }}/cron/debmirror.sh
