# roles/logging/templates/etc/init/logstash-forwarder.conf
# {{ ansible_managed }}

description "Logstash Forwarding Agent"
author "Myles Steinhauser"

setuid logstash
setgid adm

respawn
respawn limit 3 30

start on runlevel [2345]
stop on runlevel [!2345]

chdir /var/lib/logstash

{% if logging.forward.logstash_forwarder.version == '0.3.1' %}
exec /opt/logstash-forwarder/bin/logstash-forwarder.sh -config /etc/logstash-forwarder.d
{% else %}
exec /opt/logstash-forwarder/bin/logstash-forwarder -config /etc/logstash-forwarder.d
{% endif %}
