# roles/ttyspy/ttyspy-server/templates/etc/cron.daily/ttyspy_compression
#!/usr/bin/env bash

set -e

/usr/bin/python /etc/ttyspy/compression.py
