# roles/common/templates/etc/update-motd.d/cuttle-motd
#!/bin/sh
#
# {{ ansible_managed }}

cat <<End-of-message

{{ common.motd.message }}

End-of-message
